# server/com.claude-notify.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-notify.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/admin/.volta/bin/node</string>
    <string>/Users/admin/dev2/claude-notify/server/index.mjs</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/claude-notify-server.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-notify-server.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CLAUDE_NOTIFY_PORT</key>
    <string>9876</string>
  </dict>
</dict>
</plist>
